# docker-reload@.service
[Unit]
Description=Refresh docker compose container %i

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-compose@%i.service
